# archrepo.service
[Unit]
Description=archrepo server

[Service]
ExecStart=/usr/bin/archrepo -repo /home/naitik/pkgs/repo
